# luci-app-blink/root/etc/uci-defaults/22-blink
#!/bin/sh

# Ensure the init.d script is executable
chmod 777 /etc/init.d/blink

/etc/init.d/blink enable
/etc/init.d/blink start

